# systemd/telemetry-zombie-watch.service
[Unit]
Description=Telemetry zombie connection watchdog
Wants=telemetry.service

[Service]
Type=oneshot
ExecStart=/home/telemetry/scripts/telemetry-zombie-check-v2.sh
